# rootfs/etc/s6-overlay/s6-rc.d/config-maxmind/run
#!/command/with-contenv sh
if [ ! -z ${MAXMIND_KEY} ]
then
  exec s6-setuidgid tundra /scripts/maxmind_refresh.sh
fi