# rootfs/etc/s6-overlay/s6-rc.d/config-nginx/run
#!/command/with-contenv sh
###############################################################################
# Build maxmind template
###############################################################################
if [ ! -z ${MAXMIND_KEY} ]
then
  for template in 'maxmind-city.conf.template' 'maxmind-country.conf.template'
  do
    echo "[info] Enabling ${template}"
    envsubst < /templates/${template} > /etc/nginx/conf.d/${template//.template}
  done
fi

###############################################################################
# ACME
###############################################################################
ACCOUNT_THUMBPRINT=
SERVER=${ACME_SERVER:-'letsencrypt'}
if [ ! -z ${ACME_EMAIL} ]
then
  echo "[info] Configuring acme with email ${ACME_EMAIL} and server ${SERVER}"
  output=$(acme --register-account --server ${SERVER} -m ${ACME_EMAIL})

  # Extract the value of ACCOUNT_THUMBPRINT from the output using grep and awk
  ACCOUNT_THUMBPRINT=$(echo "$output" | grep "ACCOUNT_THUMBPRINT" | awk -F "'" '{print $2}')
  echo "[info] Registered! Account thumbprint is ${ACCOUNT_THUMBPRINT}"
fi

export "ACME_THUMBPRINT=${ACCOUNT_THUMBPRINT}"

echo "[info] Setting default site config"
# envsubst < /templates/default.conf.template > /etc/nginx/sites.d/0-default.conf