# rootfs/etc/s6-overlay/s6-rc.d/config-nginx/up
#!/bin/execlineb -P
# exec s6-setuidgid tundra /etc/s6-overlay/s6-rc.d/config-nginx/run
exec /etc/s6-overlay/s6-rc.d/config-nginx/run